# teslamate/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
