# teslamate/rootfs/etc/services.d/teslamate/run
#!/usr/bin/with-contenv bash
exec /bin/teslamate
